# vibeswarm.service
[Unit]
Description=VibeSwarm Web Application
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/VibeSwarm/build
ExecStart=/home/pi/VibeSwarm/build/VibeSwarm.Web
Restart=on-failure
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=vibeswarm

# Environment configuration
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://0.0.0.0:5000
Environment=HOME=/home/pi

# PATH configuration - include common CLI tool installation locations
# Modify this if your CLI tools are installed in different locations
Environment=PATH=/home/pi/.local/bin:/home/pi/.cargo/bin:/home/pi/.npm-global/bin:/home/pi/go/bin:/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin:/snap/bin

[Install]
WantedBy=multi-user.target
